# src/etc/init.d/indiestor-watch
### BEGIN INIT INFO
# Provides: indiestor-watch
# Required-Start: $all
# Required-Stop: $all
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
### END INIT INFO

case "$1" in
  start)
	echo "Start indiestor watching"
	indiestor --groups -start-watching
    ;;
  stop)
	echo "Stop indiestor watching"
	indiestor --groups -stop-watching
    ;;
  *)
    echo "Usage: /etc/init.d/indiestor-watch {start|stop}"
    exit 1
    ;;
esac

exit 0

