# src/etc/init.d/indiestor-quota
### BEGIN INIT INFO
# Provides:          indiestor-quota
# Required-Start:
# Required-Stop:
# Default-Start:
# Default-Stop:
### END INIT INFO

VOLUMES_WITH_QUOTA=/srv/indiestor.quota

case "$1" in
  start)
    echo "Starting Indiestor quota listed in $VOLUMES_WITH_QUOTA"
    if [ -f $VOLUMES_WITH_QUOTA ]; then
	while read volume; do 
		echo "Starting quota for volume $volume"
		indiestor --volume $volume -quota-on 
	done < $VOLUMES_WITH_QUOTA
    else
      	echo "no volumes with quota"
    fi
    ;;
  stop)
    echo "Stopping Indiestor quota listed in $VOLUMES_WITH_QUOTA"
    if [ -f $VOLUMES_WITH_QUOTA ]; then
	while read volume; do
		echo "Stopping quota for volume $volume"
		quotaoff --format=vfsold -u $volume 
	done < $VOLUMES_WITH_QUOTA
    else
      	echo "no volumes with quota"
    fi
    ;;
  *)
    echo "Usage: /etc/init.d/indiestor-quota {start|stop}"
    exit 1
    ;;
esac

exit 0

