# etc/init.d/indiestor-watch
### BEGIN INIT INFO
# Provides:          indiestor-quota
# Required-Start:
# Required-Stop:
# Default-Start:
# Default-Stop:
### END INIT INFO

VOLUMES_WITH_QUOTA=/srv/indiestor.quota

case "$1" in
  start)
	echo "Start indiestor watching"
	indiestor --groups -start-watching
    ;;
  stop)
	echo "Stop indiestor watching"
	indiestor --groups -stop-watching
    ;;
  *)
    echo "Usage: /etc/init.d/indiestor-watch {start|stop}"
    exit 1
    ;;
esac

exit 0

